# supervisor/hbase.conf
[program:HBase]
command=/usr/bin/start-hbase.sh
autostart=true