# supervisor/hbase.conf
[program:HBase]
command=/usr/bin/start-hbase.sh
autostart=true
autorestart=false
stderr_logfile=/opt/hbase/log/hbase.err.log
stdout_logfile=/opt/hbase/log/hbase.out.log